# services/podman/systemd/system/privileged.container
[Container]
ContainerName=nginx-example-priv
Image=docker.io/nginxinc/nginx-unprivileged
PublishPort=41080:8080
